# gofile-bot.service
[Unit]
Description=Gofile Upload Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/gofile-bot
ExecStart=/usr/bin/python3 /root/gofile-bot/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
